# pkg/systemd/jurnalo.service
[Unit]
Description=jurnalo
After=systemd-journald

[Service]
ExecStart=/usr/bin/jurnalo
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
